# deployment/supervisor/echohub-production.conf
[program:echohub-production-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/echohub/production/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --queue=default,emails,notifications
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/echohub/production/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
startsecs=0

[program:echohub-production-schedule]
process_name=%(program_name)s
command=/bin/bash -c 'while true; do php /var/www/echohub/production/artisan schedule:run --verbose --no-interaction; sleep 60; done'
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/echohub/production/storage/logs/schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
